# macos/com.ytdlp.service.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.ytdlp.service</string>

  <key>ProgramArguments</key>
  <array>
    <string>__BIN__</string>
    <string>--config</string>
    <string>__CONFIG__</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>WorkingDirectory</key>
  <string>__WORKDIR__</string>

  <key>StandardOutPath</key>
  <string>/tmp/ytdlp-service.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/ytdlp-service.error.log</string>

  <!-- Configuration is read from the config.toml path passed via ProgramArguments -->
</dict>
</plist>
